# deploy/hynous.service
[Unit]
Description=Hynous Dashboard + Agent + Daemon + Discord Bot
After=network.target nous.service
Requires=nous.service

[Service]
Type=simple
User=hynous
Group=hynous
WorkingDirectory=/opt/hynous/dashboard
EnvironmentFile=/opt/hynous/.env
Environment=PATH=/opt/hynous/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=API_URL=http://89.167.50.168:8000
ExecStart=/opt/hynous/.venv/bin/reflex run --env prod
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
